# patterns/t1/mythic/M5.temporal-sheaf/units/example.service
[Unit]
Description=Muster example for Temporal Sheaf
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern M5.temporal-sheaf'
