# patterns/t1/mythic/M5.temporal-sheaf/units/example.service
[Unit]
Description=Muster temporal sheaf obligation tracker example
Documentation=file:README.md

[Service]
Type=oneshot
User=muster
Group=muster
EnvironmentFile=-/etc/muster/temporal-sheaf.env
ExecStart=/usr/local/lib/muster/temporal-sheaf-run.sh --apply
RuntimeDirectory=muster
StateDirectory=muster/temporal-sheaf
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/muster /var/lib/muster/temporal-sheaf
