# patterns/t2/common/T2C2.dropfolder-spooler/units/example.service
[Unit]
Description=Muster example for Dropfolder Spooler
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern T2C2.dropfolder-spooler'
