# patterns/t1/common/C4.lazy-resource-gate/units/example.service
[Unit]
Description=Muster lazy resource gate readiness check
Documentation=file:README.md
After=mnt-muster-cold.automount
Wants=mnt-muster-cold.automount

[Service]
Type=oneshot
User=muster
Group=muster
ExecStart=/usr/bin/test -d /mnt/muster/cold
NoNewPrivileges=true
PrivateTmp=true
