# patterns/t1/common/C4.lazy-resource-gate/units/mnt-muster-cold.automount
[Unit]
Description=Muster lazy cold-storage automount
Documentation=file:README.md

[Automount]
Where=/mnt/muster/cold
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
