# patterns/t2/rare/T2R5.signed-update-rail/units/muster-signed-update.timer
[Unit]
Description=Run Muster signed update checks

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h
RandomizedDelaySec=20min
Persistent=true

[Install]
WantedBy=timers.target
