# patterns/t1/common/C5.failure-ratchet/units/example.service
[Unit]
Description=Muster example for Failure Ratchet
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern C5.failure-ratchet'
