# patterns/t2/rare/T2R6.home-assistant-mqtt-bridge/units/muster-ha-mqtt-bridge.timer
[Unit]
Description=Muster Home Assistant MQTT bridge timer

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
RandomizedDelaySec=5s
Persistent=true
Unit=muster-ha-mqtt-bridge.service

[Install]
WantedBy=timers.target
