# patterns/t2/common/TC3.scheduled-herald/units/example.service
[Unit]
Description=Muster scheduled herald example
Documentation=file:README.md

[Service]
Type=oneshot
User=muster
Group=muster
EnvironmentFile=-/etc/muster/herald.env
ExecStart=/usr/local/lib/muster/collect-status.sh --apply
RuntimeDirectory=muster
StateDirectory=muster/herald
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/muster /var/lib/muster
