# patterns/t2/rare/T2R4.device-triggered-conveyor/units/muster-device-conveyor@.service
[Unit]
Description=Muster device-triggered conveyor for /dev/%I
Documentation=file:README.md
After=dev-%i.device mnt-muster-cold.automount
BindsTo=dev-%i.device
Wants=mnt-muster-cold.automount
OnFailure=muster-recover@%n.service

[Service]
Type=oneshot
User=muster
Group=muster
EnvironmentFile=-/etc/muster/device-conveyor.env
ExecStart=/usr/local/lib/muster/device-convey.sh /dev/%I --apply
TimeoutStartSec=12h
RuntimeDirectory=muster
StateDirectory=muster/device-conveyor
CacheDirectory=muster/hot
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/muster /var/lib/muster /var/cache/muster /mnt/muster/cold
