# patterns/t2/common/TC4.self-healing-reconnector/units/example.service
[Unit]
Description=Muster example for Self-Healing Reconnector
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern TC4.self-healing-reconnector'
