# patterns/t1/common/C1.service-capsule/units/example.service
[Unit]
Description=Muster example for Service Capsule
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern C1.service-capsule'
