# patterns/t3/mythic/T3M2.ritualized-recovery-loop/units/example.service
[Unit]
Description=Muster example for Ritualized Recovery Loop
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern T3M2.ritualized-recovery-loop'
