# patterns/t1/rare/R2.device-binding/units/muster-device-bound@.service
[Unit]
Description=Muster device-bound worker for /dev/%I
Documentation=file:README.md
After=dev-%i.device
BindsTo=dev-%i.device
OnFailure=muster-recover@%n.service

[Service]
Type=oneshot
User=muster
Group=muster
EnvironmentFile=-/etc/muster/device-binding.env
ExecStart=/usr/local/lib/muster/device-bound-run.sh /dev/%I --apply
RuntimeDirectory=muster
StateDirectory=muster/device-binding
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/muster /var/lib/muster
