# patterns/t1/rare/R4.state-ledger/units/example.service
[Unit]
Description=Muster example for State Ledger
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern R4.state-ledger'
