# patterns/t1/mythic/M1.local-truth-sheaf/units/example.service
[Unit]
Description=Muster local truth sheaf reducer example
Documentation=file:README.md

[Service]
Type=oneshot
User=muster
Group=muster
EnvironmentFile=-/etc/muster/local-truth-sheaf.env
ExecStart=/usr/local/lib/muster/local-truth-reduce.sh --apply
RuntimeDirectory=muster
StateDirectory=muster/local-truth
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/muster /var/lib/muster/local-truth
