# patterns/t1/mythic/M1.local-truth-sheaf/units/example.service
[Unit]
Description=Muster example for Local Truth Sheaf
Documentation=file:README.md

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo muster pattern M1.local-truth-sheaf'
