# patterns/t2/common/TC5.local-sidecar-bridge/units/muster-sidecar-bridge.timer
[Unit]
Description=Muster local sidecar status bridge timer
Documentation=file:README.md

[Timer]
OnBootSec=4m
OnUnitActiveSec=5m
Persistent=true
RandomizedDelaySec=45s
Unit=muster-sidecar-bridge.service

[Install]
WantedBy=timers.target
