# patterns/t2/common/TC1.hot-cold-nas-conveyor/units/muster-nas-conveyor.service
[Unit]
Description=Muster hot-to-cold NAS conveyor
Documentation=file:README.md
After=mnt-muster-cold.automount
Wants=mnt-muster-cold.automount
OnFailure=muster-recover@%n.service

[Service]
Type=oneshot
User=muster
Group=muster
EnvironmentFile=-/etc/muster/nas-conveyor.env
ExecStart=/usr/local/lib/muster/convey --once --apply
RuntimeDirectory=muster
StateDirectory=muster/nas-conveyor
CacheDirectory=muster/hot
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/muster /var/lib/muster /var/cache/muster /mnt/muster/cold
